# Billgates/1809/etc/cron.daily/debtags
#!/bin/sh

# Update the data and index files of Debtags
if [ -x /usr/bin/debtags ]
then
	/usr/bin/debtags update 2>/dev/null || echo debtags update error
fi
